# nomad.upstart.conf
exec nomad agent -config /etc/nomad
start on runlevel [2345]
stop on runlevel [!2345]
